# redcap_docker_aio/supervisord.conf
; Runs MariaDB, Mailhog, and Apache together inside the all-in-one container.
; entrypoint.sh handles first-run DB bootstrap, then execs supervisord here.
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid

[program:mariadb]
# Durability relaxed for a disposable test DB: skip per-commit fsync and the
# doublewrite buffer so the per-spec reseed loads much faster. Each spec still
# gets a FULL reset (isolation unchanged); only crash-durability is traded away.
command=/usr/sbin/mariadbd --user=mysql --datadir=/var/lib/mysql --bind-address=0.0.0.0 --innodb-flush-log-at-trx-commit=0 --sync-binlog=0 --innodb-doublewrite=0
priority=10
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mailhog]
command=/usr/local/bin/mailhog
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:apache]
command=apache2-foreground
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
